# etc/systemd/system/timeline.service
[Unit]
Description=Scrape twitter target timelines

[Service]
User=mfa
Group=mfa
Type=oneshot
ExecStart=/usr/local/bin/mfa sync-timelines
